# Resources/linux/osqueryd.service
[Unit]
Description=Vistar MDM
After=network.service syslog.service

[Service]
TimeoutStartSec=0
EnvironmentFile=/etc/profile